# modules/dhmon/templates/alertmanager.service.erb
# Managed by Puppet
#
## !! Any changes to this file will be overwritten !!

[Unit]
Description=Alertmanager
After=syslog.target network.target

[Service]
EnvironmentFile=/etc/default/alertmanager
User=prometheus
WorkingDirectory=/opt/alertmanager
ExecStart=/opt/alertmanager/alertmanager $ARGS

[Install]
WantedBy=multi-user.target

