# modules/prometheus/templates/prometheus.service.erb
# Managed by Puppet
#
## !! Any changes to this file will be overwritten !!

[Unit]
Description=Prometheus
After=syslog.target network.target

[Service]
EnvironmentFile=/etc/default/prometheus
User=prometheus
WorkingDirectory=/opt/prometheus
ExecStart=/opt/prometheus/prometheus $ARGS
ExecReload=kill -HUP $MAINPID
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

