# modules/netbox/templates/housekeeping.service.erb
[Unit]
Description=NetBox Housekeeping Service
Documentation=https://docs.netbox.dev/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

User=netbox
Group=netbox
WorkingDirectory=/opt/netbox

ExecStart=/opt/netbox/venv/bin/python /opt/netbox/netbox/manage.py housekeeping

[Install]
WantedBy=multi-user.target
