# modules/etcd/templates/etcd.service.erb
[Unit]
Description=etcd
Documentation=https://github.com/coreos/etcd

[Service]
Type=notify
Restart=always
RestartSec=5s
TimeoutStartSec=0

ExecStart=/usr/bin/etcd --name <%= @hostname %> \
    --data-dir /var/lib/etcd \
    --listen-client-urls https://<%= @ipaddress %>:2379 \
    --advertise-client-urls https://<%= @ipaddress %>:2379 \
    --listen-peer-urls https://<%= @ipaddress %>:2380 \
    --initial-advertise-peer-urls https://<%= @ipaddress %>:2380 \
    --initial-cluster <%= @nodes.map{ |node| "#{node}=https://#{node}:2380" }.join(',') %> \
    --initial-cluster-state new \
    --ca-file /etc/etcd/ca.crt \
    --cert-file /etc/etcd/client.pem \
    --key-file /etc/etcd/key.pem \
    --client-cert-auth \
    --peer-trusted-ca-file /etc/etcd/ca.crt \
    --peer-cert-file /etc/etcd/peer.pem \
    --peer-key-file /etc/etcd/key.pem \
    --peer-client-cert-auth
 
[Install]
WantedBy=multi-user.target