# modules/supervisor/templates/include.ini.erb
# AUTOGENERATED BY PUPPET
# All manual changes will be overwritten
#
# Puppet generated configuration file
# See http://supervisord.org/configuration.html for details

[program:<%= @name %>]
command=<%= @command %>
<% if not @user.empty? -%>
user=<%= @user %>
<% end -%>
<% if not @directory.empty? -%>
directory=<%= @directory %>
<% end -%>
<% if not @autostart.to_s.empty? -%>
autostart=<%= @autostart %>
<% end -%>
<% if not @stopasgroup.to_s.empty? -%>
stopasgroup=<%= @stopasgroup %>
<% end -%>
<% if not @environment.empty? -%>
environment=<%= @environment %>
<% end -%>
