# modules/dhmoncolo/templates/prometheus-colo.service.erb
# Managed by Puppet
#
## !! Any changes to this file will be overwritten !!

[Unit]
Description=Prometheus Colo
After=syslog.target network.target

[Service]
EnvironmentFile=-/etc/defaults/prometheus-colo
User=prometheus
WorkingDirectory=/opt/prometheus-colo
ExecStart=/opt/prometheus-colo/prometheus $ARGS

[Install]
WantedBy=multi-user.target

