# modules/dhmoncolo/templates/prometheus-event.service.erb
# Managed by Puppet
#
## !! Any changes to this file will be overwritten !!

[Unit]
Description=Prometheus Colo
After=syslog.target network.target

[Service]
EnvironmentFile=-/etc/defaults/prometheus-event
User=prometheus
WorkingDirectory=/opt/prometheus-event
ExecStart=/opt/prometheus-event/prometheus $ARGS

[Install]
WantedBy=multi-user.target

