# modules/promtail/templates/promtail.service.erb
[Unit]
Description=Promtail
Documentation=https://github.com/grafana/loki

[Service]
Restart=always
User=promtail
EnvironmentFile=/etc/default/promtail
ExecStart=/opt/promtail/promtail $ARGS
ExecReload=/bin/kill -HUP $MAINPID
TimeoutStopSec=60s
SendSIGKILL=no

[Install]
WantedBy=multi-user.target
